# ansible/roles/docker/config/supervisord.conf
[program:docker]
command=docker run -d -p 8080:80 uniqcast:latest
